# step3.service
[Unit]
Before=network.target

[Service]
ExecStart=/usr/local/bin/step3.sh

[Install]
WantedBy=network.target